# crypto_portfolio/rootfs/etc/services.d/crypto_portfolio/finish
#!/usr/bin/with-contenv bashio

# This script is executed when the service stops
